# build/stratos2/puppet/modules/php_cartridge/files/etc/rc.local
#!/bin/bash -e
mkdir -p /etc/agent/conf
chmod 755 /opt/wso2-openstack-init.sh
/opt/wso2-openstack-init.sh > /var/log/wso2-openstack-init.log
exit 0
    
